# packaging/deb/clipvault_1.2.0/etc/xdg/autostart/clipvault.desktop
[Desktop Entry]
Name=ClipVault
Comment=Clipboard history daemon
Exec=clipvault
Icon=clipvault
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=3
